# vagrant/cookbooks/sniffmypackets/files/default/sniffmypackets.conf
[program:sniffmypackets-web]
command=python server.py
directory=/sniffmypackets
stdout_logfile=/home/vagrant/sniffmypackets-web.log
stderr_logfile=/home/vagrant/sniffmypackets-error.log